# mcp-server/deploy/l0g-mcp-deploy.service
[Unit]
Description=Deploy attested l0g.fr MCP release atomically
Documentation=https://github.com/bluetouff/l0g/blob/main/docs/MCP-RELEASE.md
After=network-online.target
Wants=network-online.target
ConditionPathIsDirectory=/opt/l0g-mcp-runtime

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/l0g-mcp-deploy
TimeoutStartSec=12min
UMask=0077

NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
PrivateIPC=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
ReadWritePaths=/opt/l0g-mcp-runtime
LimitNOFILE=4096
LimitCORE=0
MemoryMax=512M
TasksMax=128
